# cdn/fivem-server.service
[Unit]
Description=IgnitionHost Fivem Server

[Service]
Type=forking
ExecStart=/usr/bin/fivem_ignition.sh

[Install]
WantedBy=multi-user.target