# cdn/fivem-server.service
[Unit]
Description=NOPING Fivem Server

[Service]
Type=forking
ExecStart=/usr/bin/fivem_noping.sh

[Install]
WantedBy=multi-user.target
